# systemd/mcp-server.service
[Unit]
Description=Claude MCP File Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/mcp/app

Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/ec2-user
EnvironmentFile=/opt/mcp/app/.env

ExecStart=/opt/mcp/bin/mcp-server sse
Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-server

[Install]
WantedBy=multi-user.target
